# fanshim/rootfs/etc/services.d/fanshim/run
#!/usr/bin/with-contenv bash
exec python3 /run.py
